# deploy/qmtquant-api.service
[Unit]
Description=QmtQuant Data API 服务
Documentation=https://github.com/your-org/QmtQuant
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/data2/qmtquant

# 从 .env 文件自动加载环境变量（DATA_API_HOST、DATA_API_PORT、XQSHARE_REMOTE_HOST 等）
EnvironmentFile=/data2/qmtquant/.env

# 使用虚拟环境中的 Python 启动 data-api 服务
ExecStart=/data2/qmtquant/venv/bin/python dm_cli.py data-api

# 崩溃后自动重启，等待 10 秒
Restart=on-failure
RestartSec=10

# 日志输出到 journald（通过 journalctl -u qmtquant-api 查看）
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qmtquant-api

# 进程优先级（可选，适当降低以避免影响其他服务）
Nice=10

[Install]
WantedBy=multi-user.target
